# 우분투 자동스크립트/keyroot.service
[Unit]
Description=Key Root Community Flask App
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/keyroot_community
Environment="PATH=/var/www/keyroot_community/venv/bin"
ExecStart=/var/www/keyroot_community/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 wsgi:app
Restart=always

[Install]
WantedBy=multi-user.target